# i3battery/i3battery.service
[Unit]
Description=Battery warning for i3

[Service]
Type=oneshot
ExecStart=%h/bin/i3battery.sh

[Install]
WantedBy=default.target
