# i3battery/i3battery.timer
[Unit]
Description=Periodic check for low and unplugged battery

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
